# production/config/systemd/intent-recognition.service
[Unit]
Description=Athena Intent Recognition Service
Documentation=https://github.com/athena/intent-service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify
User=athena
Group=athena
WorkingDirectory=/opt/athena
Environment="PATH=/opt/athena/venv/bin"
EnvironmentFile=-/opt/athena/.env.production

# 服务配置
ExecStart=/opt/athena/venv/bin/gunicorn \
    core.api.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --access-logfile /var/log/athena/intent-service/access.log \
    --error-logfile /var/log/athena/intent-service/error.log \
    --log-level info \
    --timeout 120 \
    --keepalive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100

# 进程管理
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/athena/data /opt/athena/cache /var/log/athena

[Install]
WantedBy=multi-user.target
