# production/launchd/com.athena.knowledge-graph.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.athena.knowledge-graph</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/xujian/Athena工作平台/production/scripts/knowledge_graph_service.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/xujian/Athena工作平台/production/logs/knowledge_graph.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/xujian/Athena工作平台/production/logs/knowledge_graph.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/xujian/Athena工作平台</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>KG_HOST</key>
        <string>0.0.0.0</string>
        <key>KG_PORT</key>
        <string>8002</string>
    </dict>
</dict>
</plist>
